# visitors_addon/files/general/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Wait for config to be ready
bashio::log.info "Waiting for configuration..."
while [ ! -f /visitors_config/options.json ]; do
    sleep 1
done

bashio::log.info "Starting nginx..."
exec nginx
